# symlinks/config/systemd/user/wallpaper.service
[Unit]
Description=wallpaper

Wants=xsession.target
After=xsession.target

ConditionFileIsExecutable=%h/.wallpaper/wallpaper.sh

[Service]
Type=oneshot

ExecStart=%h/.wallpaper/wallpaper.sh

[Install]
WantedBy=xsession.target
