# symlinks/config/systemd/user/redshift.service.d/xsession.conf
[Unit]
Wants=xsesision.target
After=xsesision.target

[Service]
Environment="DISPLAY=:0"
RestartSec=5
